# templates/systemd/zabbix_agentd.service.j2
[Unit]
Description=
Documentation=
Wants=
After=

[Service]
Environment=
EnvironmentFile=
WorkingDirectory=
User=
Group=
ExecStartPre=
ExecStart=

[Install]
WantedBy=
